# infra/systemd/tasks-stg.service
[Unit]
Description=Tasks API (staging)
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/tasks-stg/apps/api
EnvironmentFile=/etc/tasks-stg.env
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
