# infra/systemd/tasks-stg.service
[Unit]
Description=Tasks API (staging)
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/tasks-stg/apps/api
EnvironmentFile=/etc/tasks-stg.env
Environment=NODE_ENV=production
Environment=PORT=3102
Environment=HOST=127.0.0.1
Environment=REQUEST_TIMEOUT_MS=15000
ExecStartPre=/usr/bin/env node scripts/migrate.js
ExecStart=/usr/bin/node src/server.js
ExecStartPost=/bin/sh -c 'sleep 2 && curl -fsS http://127.0.0.1:3102/api/health > /dev/null'
Restart=on-failure
RestartSec=5
KillMode=process
TimeoutStopSec=30
RuntimeDirectory=tasks-stg
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
