# infra/systemd/tasks-prod.service
[Unit]
Description=Tasks API (production)
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/tasks-prod/apps/api
EnvironmentFile=/etc/tasks-prod.env
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
